# containers/kasm-base/conf.d/lwp-mounts.conf
[program:lwp-mounts]
; Mounts the user's extra SFTP/S3 storages under ~/Mount/<name> (rclone).
; No-op (sleeps) when LWP_MOUNTS is unset.
environment=HOME="/home/lwp",USER="lwp"
command=/usr/local/bin/lwp-mounts.sh
autorestart=unexpected
exitcodes=0
startsecs=3
priority=105
user=lwp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
