# containers/kasm-base/conf.d/lwp-vpn-relay.conf
[program:lwp-vpn-relay]
; Per-session SOCKS5 switchboard (127.0.0.1:1081): DIRECT or via the user's
; VPN gateway, toggled from the window titlebar. Exits 0 when the session was
; launched without a live gateway (LWP_VPN_UPSTREAM unset).
environment=HOME="/home/lwp",USER="lwp"
command=/usr/local/bin/lwp-vpn-relay.py
autorestart=unexpected
exitcodes=0
startsecs=0
priority=230
user=lwp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
