# containers/kasm-base/conf.d/lwp-opener.conf
[program:lwp-opener]
; Polls the backend and opens "Open with…" files in the running app instance.
environment=HOME="/home/lwp",USER="lwp"
command=/usr/local/bin/lwp-opener.sh
autorestart=true
startsecs=1
priority=220
user=lwp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
