# containers/kasm-base/conf.d/lwp-audio.conf
[program:lwp-audio]
; Opus/Ogg HTTP audio stream of the desktop sink monitor on :8081.
environment=HOME="/home/lwp",USER="lwp",XDG_RUNTIME_DIR="/run/user/1000"
command=/usr/local/bin/lwp-audio.sh
autorestart=true
startsecs=1
priority=200
user=lwp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
