# containers/kasm-base/conf.d/lwp-video.conf
[program:lwp-video]
; Raw H.264 stream of the desktop on :8082 for the WebCodecs beta viewer.
; ffmpeg sits in HTTP listen (no encoding) until a client connects.
environment=HOME="/home/lwp",USER="lwp"
command=/usr/local/bin/lwp-video.sh
autorestart=true
startsecs=1
priority=210
user=lwp
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
